# systemd-unit/credfeto-batch-updates.timer
[Unit]
Description=Update Credfeto Batch Updates Repo

[Timer]
OnBootSec=5min
OnCalendar=*:0/15
AccuracySec=30

[Install]
WantedBy=multi-user.target
